# .config/systemd/user/anki.service
[Unit]
Description=Anki
After=wm.target

[Service]
.include ../config/XAppEnv
ExecStart=/usr/bin/anki
