# .config/systemd/user/compton.service
[Unit]
Description=Compton
Before=after-wm.target
After=wm.target

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=after-wm.target
