# .config/systemd/user/empathy.service
[Unit]
Description=Empathy

[Service]
ExecStart=/usr/bin/empathy

[Install]
WantedBy=wm.target
