# .config/systemd/user/prepare.target
[Unit]
Description=Prepare stage
Wants=cache-dir.service
