# .config/systemd/user/udisksvm.service
[Unit]
Description=Udisksvm
Wants=wm.target

[Service]
ExecStart=/usr/bin/udisksvm

[Install]
WantedBy=after-wm.target
