# .config/systemd/user/xcompmgr.service
[Unit]
Description=Xcompmgr
Before=after-wm.target
After=wm.target

[Service]
ExecStart=/usr/bin/xcompmgr

[Install]
WantedBy=after-wm.target
